# deploy/nexwatch-agent.service
[Unit]
Description=NexWatch Agent — system monitoring collector
Documentation=https://github.com/CogniDevAI/nexwatch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nexwatch
# Uncomment if docker group exists on the host:
# SupplementaryGroups=docker

# Config: auto-discovers /etc/nexwatch/agent.yaml
# Optional env-var overrides (NEXWATCH_* prefix). File is optional — systemd won't fail if absent.
EnvironmentFile=-/etc/nexwatch/agent.env

ExecStart=/usr/local/bin/nexwatch-agent

Restart=on-failure
RestartSec=5

# File descriptor limit (needed for high-cardinality metric collection)
LimitNOFILE=65536

SyslogIdentifier=nexwatch-agent

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Allow reading /proc and /sys for metric collectors (diskio, connections, services)
ProtectProc=invisible
ProcSubset=all
ReadOnlyPaths=/proc /sys

[Install]
WantedBy=multi-user.target
